# _uberspace/fotohaeckertwo.ini
[program:fotohaeckertwo]
command=mix phx.server
directory=%(ENV_FH_DIR)s
environment=
    MIX_ENV="prod",
    PORT=%(ENV_FH_PORT)s,
    DATABASE_PATH=%(ENV_FH_DB_PATH)s,
    SECRET_KEY_BASE=%(ENV_FH_SECRET_KEY_BASE)s,
    PHX_HOST=%(ENV_FH_PHX_HOST)s
autostart=yes
autorestart=yes